# systemd/canary.service
[Unit]
Description=Canary Server
After=mysql.service
StartLimitIntervalSec=1
[Service]
Type=simple
Restart=always
RestartSec=1

WorkingDirectory=/home/ubuntu/canary
ExecStart=/bin/bash /home/ubuntu/canary/start.sh
User=ubuntu

[Install]
WantedBy=multi-user.target
